# etc/local.d/disable_bluetooth.start
echo disable > /proc/acpi/ibm/bluetooth
sleep 1
rmmod btusb